# dist/waylaunchd.service.in
[Unit]
Description=waylaunch content-search indexer
Documentation=man:waylaunch(1)
# Index while the graphical session is up; queries work even if this is down.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Paths are substituted from the install prefix at configure time (a plain
# /usr/bin here would break a /usr/local install with status=203/EXEC).
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/waylaunchd
ExecReload=@CMAKE_INSTALL_FULL_BINDIR@/waylaunchctl reconcile
Restart=on-failure
RestartSec=5

# Resource governance (docs/CONTENT_SEARCH.md §7 / NFR3/NFR4): run indexing at
# background priority and cap memory so it is safe to leave running on a laptop.
Nice=10
IOSchedulingClass=idle
CPUSchedulingPolicy=batch
# Delegate a cgroup subtree so the daemon can sandbox extractor subprocesses in
# a child cgroup (memory.max / pids.max / oom.group + clone3 CLONE_INTO_CGROUP).
# Requires write access to the subtree, so ProtectControlGroups must stay off.
Delegate=cpu memory pids
# Envelope for the whole tree (daemon + extractors). The daemon confines the
# extractors further via its ./extract subgroup (default 512M, swap off); no
# MemoryHigh here — throttling a PDF parser mid-file just makes it hang longer.
MemoryMax=768M
# Hardening: the daemon needs no elevated privileges and no network.
PrivateNetwork=true
NoNewPrivileges=true
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX
ReadWritePaths=%h/.local/share/waylaunch %t/waylaunch

[Install]
WantedBy=graphical-session.target
